# hass-addon-control-group/rootfs/etc/s6-overlay/s6-rc.d/cgroup/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting the Control Group add-on"
uv run --no-sync python -m ha_addon_cgroup
