# hass-addon-esp/rootfs/etc/s6-overlay/s6-rc.d/esp/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting the ESP add-on"
uv run --no-sync python -m ha_addon_esp
