# hass-addon-qsusb64/rootfs/etc/s6-overlay/s6-rc.d/qsusb64/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting the QSUSB64 add-on"
uv run --no-sync python -m ha_addon_qsusb64
